# kramden/22.04/scripts/rc.local
#!/bin/bash

/home/snaps.sh

if [ $? -eq 0 ]; then
	rm -rf /home/snap*
	rm /etc/rc.local
fi

# Restart networking to pickup interface changes
systemctl restart NetworkManager

# Disable cloud-init as it blocks boot until timeout
touch /etc/cloud/cloud-init.disabled
