# kramden/22.04/scripts/rc.local
#!/bin/bash

/home/snaps.sh

if [ $? -eq 0 ]; then
	rm -rf /home/snap*
	rm /etc/rc.local
	nmcli n off
	sleep 2
	nmcli n on
fi
